# config/com.larcity.cami.plist.erb
<%
require 'lib/app_utils'
%>
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.larcity.cami</string>
    <key>ProgramArguments</key>
    <array>
      <string><%= Rails.root.join("bin/start").to_s %></string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string><%= Rails.root.join("log/#{ENV['RAILS_ENV']}.log").to_s %></string>
    <key>StandardErrorPath</key>
    <string><%= Rails.root.join("log/#{ENV['RAILS_ENV']}.log").to_s %></string>
</dict>
</plist>
